# rc.local
#! /bin/bash
/etc/scripts/addarp.sh
systemctl start firewall.service
exit 0
