# rc.local
#! /bin/bash
/etc/scripts/addarp.sh
exit 0
